# data_assistant_project/new_data_assistant_project/deployment/docker-systemd.service
[Unit]
Description=Data Assistant Docker Containers
Requires=docker.service
After=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/%i/bachloreArbeit/data_assistant_project/new_data_assistant_project/deployment
ExecStart=/usr/local/bin/docker-compose up -d --profile production
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=300
TimeoutStopSec=120

# Restart-Politik
Restart=on-failure
RestartSec=30

# Umgebung
Environment=COMPOSE_PROJECT_NAME=data_assistant

[Install]
WantedBy=multi-user.target 